# docker/supervisord.conf
cat > docker/supervisord.conf <<'EOF'
[supervisord]
nodaemon=true

[program:string_search_api]
directory=/app
command=/opt/conda/envs/autodfbench/bin/python -m API.string_search_api
autorestart=true
startretries=10

[program:deleted_file_recovery_api]
directory=/app
command=/opt/conda/envs/autodfbench/bin/python -m API.deleted_file_recovery_api
autorestart=true
startretries=10

[program:file_carving_api]
directory=/app
command=/opt/conda/envs/autodfbench/bin/python -m API.file_carving_api
autorestart=true
startretries=10

[program:windows_registry_api]
directory=/app
command=/opt/conda/envs/autodfbench/bin/python -m API.windows_registry_api
autorestart=true
startretries=10

[program:sqlite_recovery_api]
directory=/app
command=/opt/conda/envs/autodfbench/bin/python -m API.sqlite_recovery_api
autorestart=true
startretries=10
EOF
